# server/supervisor.conf
[program:h404bi.com.server]
directory=/data/backend/h404bi.com/server
command=sh -c './.venv/bin/gunicorn manager:entry -c ./gunicorn_cfg.py'
stopasgroup=true
user=www
